# crates/service/debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload
    systemctl enable kmitlnetauth.service
    systemctl start kmitlnetauth.service || true
fi
